# supervisord/gearman-consumer.ini
[program:updateProfile]
command=php56 /opt/www/cli.php GearmanWeixin updateProfile local
logfile=/var/log/supervisor/updateProfile.log